# gnote.go/misc/systemd/gnote.user.service
[Unit]
Description=Gnote Service
After=local-fs.target

[Service]
ExecStart=/usr/bin/gnote watch

[Install]
WantedBy=default.target
